# deploy/systemd/task-aversion-app.service
# systemd service for Task Aversion System
# Install: sudo cp deploy/systemd/task-aversion-app.service /etc/systemd/system/
# Then edit paths and: sudo systemctl daemon-reload && sudo systemctl enable task-aversion-app
#
# Replace APP_DIR with your actual path, e.g.:
#   /opt/task-aversion-system/task_aversion_app
#   /home/your-user/task_aversion_app

[Unit]
Description=Task Aversion System
After=network.target postgresql.service

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=APP_DIR
Environment="PATH=APP_DIR/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=APP_DIR/.env.production
ExecStart=APP_DIR/venv/bin/python app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
